# packaging/providers/luks/systemd/lockchain-luks.service
[Unit]
Description=LockChain Daemon (legacy unit name: lockchain-luks.service)
After=run-lockchain.mount network-online.target
Wants=network-online.target run-lockchain.mount
Requires=run-lockchain.mount
Conflicts=lockchain.service lockchain-zfs.service
ConditionPathExists=/etc/lockchain.toml

[Service]
Type=simple
User=lockchain
Group=lockchain
WorkingDirectory=/var/lib/lockchain
RuntimeDirectory=lockchain
RuntimeDirectoryMode=0750
Environment=LOCKCHAIN_CONFIG=/etc/lockchain.toml
Environment=LOCKCHAIN_HEALTH_ADDR=127.0.0.1:8787
ExecStart=/usr/bin/lockchain-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
